# deploy/supervisord.dir/dungeons.ini
[program:dungeons]  
directory=../api/
command=python server.py
autostart=false  
autorestart=true
priority=110
redirect_stderr=true

